# scripts/crontab
0 1 * * * /bin/bash /usr/src/app/scripts/refresh_data.sh >> /var/log/cron.log
