# scripts/crontab
* * 1 * * ./usr/src/app/scripts/refresh_data.sh >> /var/log/cron.log
