# deployment_scripts/crontab
* * 1 * * ./usr/src/app/refresh_data.sh >> /var/log/cron.log
